# ansible/roles/hindsight/templates/hindsight-proxy.service.j2
{# Managed by Ansible (claude-devbox). One socat proxy per profile, exposing that
   profile's loopback Hindsight daemon to the TAILNET only (UFW scopes it to
   tailscale0). Listens on the external port, forwards to the profile's loopback port. #}
[Unit]
Description=Hindsight memory API proxy ({{ item.user }}) — tailnet only
After=network-online.target tailscaled.service
Wants=network-online.target

[Service]
ExecStart=/usr/bin/socat TCP-LISTEN:{{ (hindsight_expose_base_port | default(19077)) + idx }},fork,reuseaddr TCP:127.0.0.1:{{ (hindsight_base_port | default(9077)) + idx }}
Restart=always
RestartSec=2
User=nobody

[Install]
WantedBy=multi-user.target
