# ansible/roles/hindsight/templates/hindsight-control-plane.service.j2
{# Managed by Ansible (claude-devbox). Per-profile Hindsight Control Plane web UI
   (npx @vectorize-io/hindsight-control-plane) fronting that profile's loopback
   daemon. Bound to the tailnet (0.0.0.0; public is blocked by UFW's default-deny),
   reachable at http://<node>.<tailnet>.ts.net:<ui_port>. Served DIRECTLY (not via
   Tailscale Serve) because the control-plane app mishandles a TLS-terminating proxy
   (it honors X-Forwarded-Proto: https and then TLS-proxies its own plaintext port).
   node/npx resolve via mise from the login shell. #}
[Unit]
Description=Hindsight Control Plane UI ({{ item.user }}) — tailnet only
After=network-online.target tailscaled.service
Wants=network-online.target

[Service]
User={{ item.user }}
Environment=HOME=/home/{{ item.user }}
ExecStart=/bin/bash -lc "exec npx -y @vectorize-io/hindsight-control-plane --api-url http://127.0.0.1:{{ (hindsight_base_port | default(9077)) + idx }} --port {{ (hindsight_control_plane_base_port | default(19077)) + idx }} --hostname 0.0.0.0"
Restart=always
RestartSec=5

[Install]
WantedBy=multi-user.target
